# systemd-user/omarchy-theme-light.timer
[Unit]
Description=Switch to Light theme at 07:00

[Timer]
OnCalendar=*-*-* 07:00:00
Persistent=true

[Install]
WantedBy=timers.target
